# etc/systemd/system/backup-system-flatpak-names.timer
[Unit]
Description=Backup Non-System Flatpak App Names

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=default.target

